# Advance/pxpipe/pxpipe.service.template
[Unit]
Description=pxpipe — token-metering/optical-compression proxy on 127.0.0.1:@PORT@
After=network.target

[Service]
ExecStart=@NPX@ -y pxpipe-proxy
Environment=PATH=@NODE_BIN_DIR@:/usr/bin:/bin
Environment=PORT=@PORT@
# Explicit imaging allowlist — NEVER rely on the default (it has drifted across versions, and the
# matcher family-matches suffixes: a bare `gpt-5.6` entry would image `gpt-5.6-sol`, which is a
# measured double-loss: 0/15 verbatim recall + ~32% HIGHER token cost on the OpenAI lane).
# Fable is the only model with strong measured imaged-text recall. Metering works for ALL models
# regardless of this list — non-listed traffic passes through untouched but fully ledgered.
Environment=PXPIPE_MODELS=claude-fable-5
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
